# docs/launchd/com.ea.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  EA poll agent — runs one poll cycle every 5 minutes.

  SETUP (one time):
    1. Copy this file to ~/Library/LaunchAgents/
         cp docs/launchd/com.ea.poll.plist ~/Library/LaunchAgents/

    2. Replace the three CHANGEME values:
         CHANGEME_VENV_PYTHON       absolute path, e.g. /Users/you/projects/ea/.venv/bin/python
         CHANGEME_PROJECT_DIR       absolute path, e.g. /Users/you/projects/ea
         CHANGEME_ANTHROPIC_API_KEY your sk-ant-... key

    3. Load and start:
         launchctl load ~/Library/LaunchAgents/com.ea.poll.plist
         launchctl start com.ea.poll

  USEFUL COMMANDS:
    launchctl list com.ea.poll                              # check status / last exit code
    launchctl stop com.ea.poll                              # stop immediately
    launchctl unload ~/Library/LaunchAgents/com.ea.poll.plist   # remove

  LOGS:
    Structured JSON lines  -> PROJECT_DIR/ea.log
    launchd stderr (startup errors) -> PROJECT_DIR/ea-launchd.err

  NOTE ON API KEY:
    Embedding the key in the plist stores it in plain text in ~/Library/LaunchAgents/.
    That directory is readable only by your user account, which is acceptable for a
    personal machine.  If you prefer not to store it there, use the install script
    (scripts/install-launchd.sh) — if ~/.ea-env exists it will automatically use a
    wrapper script that sources the key instead of embedding it here.
-->
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.ea.poll</string>

  <!-- Run one poll cycle every 300 seconds (5 minutes).
       Keep this in sync with poll_interval_seconds in config.toml. -->
  <key>StartInterval</key>
  <integer>300</integer>

  <key>ProgramArguments</key>
  <array>
    <string>CHANGEME_VENV_PYTHON</string>
    <string>CHANGEME_PROJECT_DIR/ea.py</string>
    <string>poll</string>
    <string>--quiet</string>
  </array>

  <!-- Must match the project root so config.toml / state.json / token.json are found. -->
  <key>WorkingDirectory</key>
  <string>CHANGEME_PROJECT_DIR</string>

  <!-- stdout suppressed by --quiet; stderr captures Python startup errors. -->
  <key>StandardOutPath</key>
  <string>/dev/null</string>
  <key>StandardErrorPath</key>
  <string>CHANGEME_PROJECT_DIR/ea-launchd.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>ANTHROPIC_API_KEY</key>
    <string>CHANGEME_ANTHROPIC_API_KEY</string>
  </dict>

  <!-- Start on login and after every reboot. -->
  <key>RunAtLoad</key>
  <true/>

</dict>
</plist>
